# container-images/base/fractal-entrypoint.service
[Unit]
Description=Fractal entrypoint service, runs as root

[Service]
Type=oneshot
ExecStart=/usr/share/fractal/entry.sh

[Install]
WantedBy=multi-user.target
